# src/sysDetector-ebpf/sysDetector-monitor.service
[Unit]
Description=Process Monitor via eBPF
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/proc
Restart=always
RestartSec=5s
User=root
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target